# deploy/polymarket-crt-sweet-spot.service
[Unit]
Description=polymarket-crt-tools sweet-spot paper watcher (JSONL, no orders)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/polymarket-crt-tools
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/polymarket-crt-tools/.venv/bin/python scripts/watch_sweet_spot.py --asset btc --tf 15 --price-source pyth --interval-sec 5 --journal var/watch_sweet_spot.jsonl --crt-preset loose --crt-no-htf-filter --crt-sweep-conflict prefer_bull --sticky-arm
Restart=always
RestartSec=20

[Install]
WantedBy=multi-user.target
